# services/cron/config/crontab
*/15 * * * * echo HostingSignal cron heartbeat >> /var/log/cron.log
