# systemd/hostingsignal-daemon.service
[Unit]
Description=HostingSignal Background Daemon (Task Queue Consumer)
Documentation=https://docs.hostingsignal.com
After=network.target postgresql.service redis.service rabbitmq-server.service
Wants=redis.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/usr/local/hostingsignal/backend
Environment="PATH=/usr/local/hostingsignal/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/usr/local/hostingsignal/backend"
ExecStart=/usr/local/hostingsignal/backend/venv/bin/python -m app.daemon.scheduler
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hostingsignal-daemon

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
