# systemd/hostingsignal-web.service
[Unit]
Description=HostingSignal Next.js Frontend Web Server
Documentation=https://docs.hostingsignal.com
After=network.target hostingsignal-api.service
Wants=hostingsignal-api.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/usr/local/hspanel/developer-panel/web
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="HOSTNAME=0.0.0.0"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/node /usr/local/hspanel/developer-panel/web/.next/standalone/server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hostingsignal-web

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
