# systemd/hostingsignal-daemon.service
[Unit]
Description=HostingSignal Background Daemon (Task Queue Consumer)
Documentation=https://docs.hostingsignal.com
After=network.target hostingsignal-api.service
Wants=hostingsignal-api.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/usr/local/hspanel
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/perl /usr/local/hspanel/daemon/hs-srvd.pl
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hostingsignal-daemon

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
