# systemd/hostingsignal-monitor.service
[Unit]
Description=HostingSignal System Monitor Service
Documentation=https://docs.hostingsignal.com
After=network.target redis.service
Wants=redis.service

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/usr/local/hspanel
Environment="PATH=/usr/local/hspanel/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/local/hspanel/backend/venv/bin/python /usr/local/hspanel/append_system_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hostingsignal-monitor

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
